# scripts/crontab
*/19 * * * * sudo /usr/bin/ntpdate ntp.api.bz; sudo hwclock --systohc;
